# packages/acpid/acpid.service
[Unit]
Description=ACPI event daemon

[Service]
Type=forking
ExecStart=/usr/sbin/acpid -c /usr/lib/acpid/events
Restart=always
RestartSec=5

[Install]
WantedBy=preconfigured.target
