# packages/ecs-agent/etc-systemd-system-ecs.service.d.mount
[Unit]
Description=ECS agent drop-ins Directory (/etc/systemd/system/ecs.service.d)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
After=selinux-policy-files.service
Wants=selinux-policy-files.service

[Mount]
What=tmpfs
Where=/etc/systemd/system/ecs.service.d
Type=tmpfs
Options=nosuid,nodev,noexec,noatime,mode=0750,context=system_u:object_r:secret_t:s0

[Install]
WantedBy=preconfigured.target
