# packages/nvidia-k8s-device-plugin/nvidia-k8s-device-plugin.service
[Unit]
Description=Start NVIDIA kubernetes device plugin
RefuseManualStart=true
RefuseManualStop=true
After=kubelet.service
BindsTo=kubelet.service

[Service]
ExecStart=/usr/bin/nvidia-device-plugin
Type=simple
TimeoutSec=0
RestartSec=2
Restart=always
StandardError=journal+console

[Install]
WantedBy=multi-user.target
