# packages/nvidia-k8s-device-plugin/nvidia-k8s-device-plugin.service
[Unit]
Description=Start NVIDIA kubernetes device plugin
After=kubelet.service
Wants=kubelet.service

[Service]
ExecStart=/usr/bin/nvidia-device-plugin --config-file=/etc/nvidia-k8s-device-plugin/settings.yaml
Type=simple
TimeoutSec=0
RestartSec=2
Restart=always
StandardError=journal+console

[Install]
WantedBy=multi-user.target
