# packages/api/apiserver.service
[Unit]
Description=Thar API server
After=network.target storewolf.service migrator.service
Requires=network.target storewolf.service migrator.service

[Service]
Type=simple
ExecStart=/usr/bin/apiserver --datastore-path /var/lib/thar/datastore/current -v -v -v

[Install]
WantedBy=multi-user.target
