# packages/api/apiserver.service
[Unit]
Description=Thar API server
After=network.target
After=storewolf.service
Requires=network.target
Requires=storewolf.service

[Service]
Type=simple
ExecStart=/usr/bin/apiserver --datastore-path /var/lib/thar/datastore/v1 -v -v -v

[Install]
WantedBy=multi-user.target
