# packages/api/moondog.service
[Unit]
Description=Thar userdata configuration system
After=network-online.target
After=apiserver.service
Requires=network-online.target
Requires=apiserver.service

[Service]
Type=oneshot
ExecStart=/usr/bin/moondog
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
