# packages/api/storewolf.service
[Unit]
Description=Datastore creator

[Service]
Type=oneshot
ExecStart=/usr/bin/storewolf --data-store-base-path /var/lib/thar/datastore
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
