# packages/open-vm-tools/vmtoolsd.service
[Unit]
Description=VMware Tools service
Documentation=https://github.com/vmware/open-vm-tools
ConditionVirtualization=vmware

[Service]
ExecStart=/usr/bin/vmtoolsd
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
