# packages/chrony/chronyd.service
[Unit]
Description=A versatile implementation of the Network Time Protocol
Documentation=https://chrony.tuxfamily.org
After=network-online.target preconfigured.target
Wants=network-online.target preconfigured.target

[Service]
Type=simple
ExecStart=/usr/sbin/chronyd -d -F -1

[Install]
WantedBy=configured.target
