# packages/ssm/amazon-ssm-agent.service
[Unit]
Description=amazon-ssm-agent
After=network-online.target
StartLimitBurst=10
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=/
ExecStartPre=/usr/bin/docker build -t ssm /var/lib/amazon/ssm/
ExecStart=/usr/bin/docker run --net=host -v "/var/lib/thar/api.sock:/var/lib/thar/api.sock" ssm /usr/bin/amazon-ssm-agent
KillMode=process
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
