# packages/wicked/wickedd-dhcp4.service
[Unit]
Description=wicked DHCPv4 supplicant service
Requisite=dbus.service
After=local-fs.target dbus.service network-pre.target
Before=wickedd.service wicked.service network.target
PartOf=wickedd.service

[Service]
Type=notify
LimitCORE=infinity
ExecStart=/usr/libexec/wicked/bin/wickedd-dhcp4 --systemd --foreground --log-level notice
StandardError=null
Restart=on-abort
