# packages/containerd/containerd.service
[Unit]
Description=containerd container runtime
Documentation=https://containerd.io
After=network-online.target configured.target
Wants=network-online.target configured.target

[Service]
ExecStart=/usr/bin/containerd
Delegate=yes
KillMode=mixed
TimeoutSec=0
RestartSec=2
Restart=always
LimitNPROC=infinity
LimitCORE=infinity
LimitNOFILE=infinity
TasksMax=infinity

[Install]
WantedBy=multi-user.target
