# packages/cloudwatch-agent/cloudwatch-agent.conf
# Copyright Amazon.com, Inc. or its affiliates. All Rights Reserved.
# SPDX-License-Identifier: MIT

description     "Amazon CloudWatch Agent"
author          "Amazon.com"

start on (runlevel [345] and started network)
stop on (runlevel [!345] or stopping network)

normal exit 0

respawn

exec /usr/bin/start-amazon-cloudwatch-agent
post-stop exec sleep 1

