# packages/cloudwatch-agent/cloudwatch-agent.service
[Unit]
Description=Amazon cloudwatch-agent
After=network-online.target preconfigured.target
Wants=network-online.target preconfigured.target

[Service]
ExecStart=/opt/aws/amazon-cloudwatch-agent/bin/start-amazon-cloudwatch-agent
Restart=always
RestartSec=600s

[Install]
WantedBy=configured.target