# packages/release/activate-preconfigured.service
[Unit]
Description=Isolates preconfigured.target
DefaultDependencies=no
After=fipscheck.target
Requires=fipscheck.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl set-default preconfigured
ExecStart=/usr/bin/systemctl isolate default --no-block
RemainAfterExit=true
StandardOutput=tty
StandardError=inherit

[Install]
WantedBy=fipscheck.target
