# packages/release/activate-multi-user.service
[Unit]
Description=Isolates multi-user.target
After=configured.target reboot-if-required.service
Requires=configured.target reboot-if-required.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl set-default multi-user
ExecStart=/usr/bin/systemctl isolate default --no-block
RemainAfterExit=true
StandardError=journal+console

[Install]
WantedBy=configured.target
