# packages/release/configured.target
[Unit]
Description=Bottlerocket final configuration complete
After=preconfigured.target
Requires=preconfigured.target
AllowIsolate=yes

[Install]
RequiredBy=multi-user.target
