# packages/release/configured.target
[Unit]
Description=Bottlerocket user and dynamic configuration complete
After=settings-applier.service
Requires=settings-applier.service
AllowIsolate=yes
