# packages/release/configured.target
[Unit]
Description=Thar be configuration - user and dynamic configuration complete
After=moondog.service sundog.service
Requires=moondog.service sundog.service
AllowIsolate=yes
