# packages/release/etc-cni.mount
[Unit]
Description=CNI Configuration Directory (/etc/cni)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=tmpfs
Where=/etc/cni
Type=tmpfs
Options=nosuid,nodev,noexec,noatime,mode=0755

[Install]
WantedBy=preconfigured.target
