# packages/release/fipscheck.target
[Unit]
Description=FIPS integrity check
AllowIsolate=false
RefuseManualStart=true
RefuseManualStop=true

[Install]
WantedBy=preconfigured.target
