# packages/release/lib-modules.mount.in
[Unit]
Description=Kernel Modules (Read-Write)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
Wants=prepare-var.service
After=prepare-var.service
RequiresMountsFor=/var

[Mount]
What=overlay
Where=PREFIX/lib/modules
Type=overlay
Options=nosuid,nodev,noexec,noatime,lowerdir=/lib/modules,upperdir=/var/lib/kernel-modules/.overlay/upper,workdir=/var/lib/kernel-modules/.overlay/work,context=system_u:object_r:state_t:s0

[Install]
WantedBy=local-fs.target
