# packages/release/local.mount
[Unit]
Description=Local Directory (/local)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=/dev/disk/by-partlabel/BOTTLEROCKET-DATA
Where=/local
Type=ext4
# "noexec" omitted to allow containers and migrations to run
Options=defaults,nosuid,nodev,noatime,private

[Install]
WantedBy=preconfigured.target
