# packages/release/local.mount
[Unit]
Description=Local Directory (/local)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=/dev/disk/by-partlabel/THAR-DATA
Where=/local
Type=ext4
Options=defaults,noatime,nosuid,nodev

[Install]
WantedBy=local-fs.target
