# packages/release/mnt.mount
[Unit]
Description=Mnt Directory (/mnt)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
RequiresMountsFor=/local

[Mount]
What=/local/mnt
Where=/mnt
Type=none
Options=rbind,rshared

[Install]
WantedBy=local-fs.target
