# packages/release/opt-cni-bin.mount
[Unit]
Description=CNI Plugin Directory (/opt/cni/bin)
Conflicts=umount.target
RequiresMountsFor=/var
Before=local-fs.target umount.target

[Mount]
What=overlay
Where=/opt/cni/bin
Type=overlay
# "noexec" omitted because containerd needs to execute CNI plugins
Options=nosuid,nodev,noatime,lowerdir=/usr/libexec/cni/bin,upperdir=/var/lib/cni-plugins/.overlay/upper,workdir=/var/lib/cni-plugins/.overlay/work,context=system_u:object_r:local_t:s0

[Install]
WantedBy=preconfigured.target
