# packages/release/opt.mount
[Unit]
Description=Opt Directory (/opt)
DefaultDependencies=no
Conflicts=umount.target
Wants=prepare-opt.service mask-local-opt.service
Before=local-fs.target umount.target prepare-opt.service mask-local-opt.service
RequiresMountsFor=/local

[Mount]
What=/local/opt
Where=/opt
Type=none
Options=rbind,rshared

[Install]
WantedBy=local-fs.target
