# packages/release/opt.mount
[Unit]
Description=Opt Directory (/opt)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
Wants=prepare-local.service
After=prepare-local.service

[Mount]
What=/local/opt
Where=/opt
Type=none
Options=bind

[Install]
WantedBy=local-fs.target
