# packages/release/prepare-var-lib-thar.service
[Unit]
Description=Prepare Private Directory (/var/lib/thar)
RequiresMountsFor=/var
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=-/usr/bin/mkdir -p /var/lib/thar
RemainAfterExit=false
