# packages/release/usr-bin.mount.in
[Unit]
Description=Binaries
DefaultDependencies=no
Conflicts=umount.target
Before=umount.target

[Mount]
What=overlay
Where=PREFIX/bin
Type=overlay
Options=lowerdir=PREFIX/fips/bin:PREFIX/bin

[Install]
WantedBy=local-fs-pre.target
