# packages/release/usr-share-licenses.mount.in
[Unit]
Description=License files
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=PREFIX/share/bottlerocket/licenses.squashfs
Where=PREFIX/share/licenses
Type=squashfs
Options=defaults,ro,loop,noatime,nosuid,nodev,context=system_u:object_r:os_t:s0

[Install]
WantedBy=local-fs.target
