# packages/release/usr-share-licenses.mount.in
[Unit]
Description=License files
DefaultDependencies=no
Conflicts=umount.target
Before=umount.target
After=local-fs.target

[Mount]
What=PREFIX/share/bottlerocket/licenses.squashfs
Where=PREFIX/share/licenses
Type=squashfs
Options=defaults,ro,loop,nosuid,nodev,noexec,noatime,context=system_u:object_r:os_t:s0

[Install]
WantedBy=preconfigured.target
