# packages/release/usr-src-kernels.mount.in
[Unit]
Description=Kernel Development Sources (Read-Write)
Conflicts=umount.target
RequiresMountsFor=/var/lib/kernel-devel/lower
Before=local-fs.target umount.target

[Mount]
What=overlay
Where=PREFIX/src/kernels
Type=overlay
Options=noatime,nosuid,nodev,lowerdir=/var/lib/kernel-devel/lower,upperdir=/var/lib/kernel-devel/upper,workdir=/var/lib/kernel-devel/work,context=system_u:object_r:local_t:s0

[Install]
WantedBy=local-fs.target
