# packages/release/usr-src-kernels.mount.in
[Unit]
Description=Kernel Development Sources (Read-Write)
DefaultDependencies=no
Conflicts=umount.target
RequiresMountsFor=/var/lib/kernel-devel/.overlay/lower
Before=umount.target
After=local-fs.target

[Mount]
What=overlay
Where=PREFIX/src/kernels
Type=overlay
# "noexec" omitted to allow objtool to run
Options=nosuid,nodev,noatime,lowerdir=/var/lib/kernel-devel/.overlay/lower,upperdir=/var/lib/kernel-devel/.overlay/upper,workdir=/var/lib/kernel-devel/.overlay/work,context=system_u:object_r:state_t:s0

[Install]
WantedBy=preconfigured.target
