# packages/release/usr-src-kernels.mount.in
[Unit]
Description=Kernel Development Sources
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=PREFIX/share/bottlerocket/kernel-devel.squashfs
Where=PREFIX/src/kernels
Type=squashfs
Options=defaults,ro,loop

[Install]
WantedBy=local-fs.target
