# packages/release/var-lib-bottlerocket.mount
[Unit]
Description=Private Directory (/var/lib/bottlerocket)
DefaultDependencies=no
Conflicts=umount.target
RequiresMountsFor=/var
Before=local-fs.target umount.target

[Mount]
What=/dev/disk/by-partlabel/BOTTLEROCKET-PRIVATE
Where=/var/lib/bottlerocket
Type=ext4
Options=defaults,nosuid,nodev,noexec,noatime,private,context=system_u:object_r:private_t:s0

[Install]
WantedBy=preconfigured.target
