# packages/release/var.mount
[Unit]
Description=Var Directory (/var)
DefaultDependencies=no
Conflicts=umount.target
Wants=prepare-var.service mask-local-var.service
Before=local-fs.target umount.target prepare-var.service mask-local-var.service
RequiresMountsFor=/local

[Mount]
What=/local/var
Where=/var
Type=none
Options=rbind,rshared

[Install]
WantedBy=local-fs.target
