# packages/release/var.mount
[Unit]
Description=Var Directory (/var)
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
Wants=prepare-local.service
After=prepare-local.service

[Mount]
What=/local/var
Where=/var
Type=none
Options=bind

[Install]
WantedBy=local-fs.target
