# packages/workspaces/apiserver.service
[Unit]
Description=Thar API server
After=storewolf.service migrator.service
Requires=storewolf.service migrator.service

[Service]
Type=notify
ExecStart=/usr/bin/apiserver --datastore-path /var/lib/thar/datastore/current --socket-gid 274 -v -v -v

[Install]
WantedBy=multi-user.target
