# packages/workspaces/mark-successful-boot.service
[Unit]
Description=Call signpost to mark the boot as successful after all required targets are met.
After=multi-user.target
# Each service that must start correctly in order for a boot to be successful should be of type "notify"
# and include "RequiredBy=mark-successful-boot.service" in its [Install] section.

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/signpost mark-successful-boot

[Install]
WantedBy=multi-user.target
