# packages/workspaces/moondog.service
[Unit]
Description=Thar userdata configuration system
# Need network online to talk to IMDS.
After=network-online.target apiserver.service
Requires=network-online.target apiserver.service
# We only want to run once, at first boot.  This file is created by moondog
# after a successful run.
ConditionPathExists=!/var/lib/thar/moondog.ran

[Service]
Type=oneshot
ExecStart=/usr/bin/moondog
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
