# packages/workspaces/settings-committer.service
[Unit]
Description=Commit pending settings
After=storewolf.service apiserver.service sundog.service moondog.service
Requires=apiserver.service
Wants=storewolf.service sundog.service moondog.service

[Service]
Type=oneshot
ExecStart=/usr/bin/settings-committer
RemainAfterExit=true
StandardError=journal+console

[Install]
WantedBy=multi-user.target
