# packages/workspaces/storewolf.service
[Unit]
Description=Datastore creator
After=migrator.service
Requires=migrator.service

[Service]
Type=oneshot
ExecStart=/usr/bin/storewolf --data-store-base-path /var/lib/thar/datastore
RemainAfterExit=true
StandardError=journal+console

[Install]
WantedBy=multi-user.target
