# packages/workspaces/updog.service
[Unit]
Description=What's Updog?

[Service]
Type=oneshot
ExecStart=/usr/bin/updog check-update
