# packages/soci-snapshotter/soci-snapshotter.service
[Unit]
Description=soci snapshotter containerd plugin
Documentation=https://github.com/awslabs/soci-snapshotter
After=configured.target
Wants=configured.target
Before=containerd.service

[Service]
Type=notify
ExecStart=/usr/bin/soci-snapshotter-grpc
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
RequiredBy=containerd.service
