# packages/os/apiserver.service
[Unit]
Description=Bottlerocket API server
After=storewolf.service migrator.service
Requires=storewolf.service migrator.service

[Service]
Type=notify
ExecStart=/usr/bin/apiserver --datastore-path /var/lib/bottlerocket/datastore/current --socket-gid 274
Restart=always
RestartSec=5
StandardError=journal+console

[Install]
WantedBy=preconfigured.target
