# packages/os/has-boot-ever-succeeded.service
[Unit]
Description=Checks and marks if boot has ever succeeded before
DefaultDependencies=no
Before=label-data-a.service label-data-b.service
RequiresMountsFor=/etc

[Service]
Type=oneshot
# Check if boot has ever succeeded before
ExecStart=/bin/signpost has-boot-ever-succeeded
RemainAfterExit=true
# If boot has succeeded before, the marker file will be non-zero-sized
StandardOutput=file:/etc/has-boot-ever-succeeded

[Install]
WantedBy=local-fs.target
