# packages/os/metricdog.service
[Unit]
Description=Send a Metricdog Ping

[Service]
EnvironmentFile=/etc/network/proxy.env
Type=oneshot
RemainAfterExit=false
StandardError=journal+console
ExecStart=/usr/bin/metricdog send-health-ping
TimeoutStartSec=30s
