# packages/os/migrator.service
[Unit]
Description=Bottlerocket data store migrator

[Service]
Type=oneshot
ExecStart=/usr/bin/migrator --datastore-path /var/lib/bottlerocket/datastore/current --migration-directories /var/lib/bottlerocket-migrations --migrate-to-version-from-os-release
RemainAfterExit=true
StandardError=journal+console

[Install]
WantedBy=multi-user.target
